# deploy/captcha-service.service
[Unit]
Description=VK captcha solver / TURN credential proxy
Documentation=https://github.com/cacggghp/vk-turn-proxy/tree/main/captcha-service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/vk-turn-proxy/captcha-service.env
ExecStart=/opt/vk-turn-proxy/captcha-service

Restart=always
RestartSec=5s
StartLimitIntervalSec=0

DynamicUser=yes
StateDirectory=vk-turn-proxy-captcha
Environment=CAPTCHA_TRAP_DIR=/var/lib/vk-turn-proxy-captcha/trap
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
PrivateTmp=true
PrivateDevices=true
RestrictSUIDSGID=true
RestrictRealtime=true
RestrictNamespaces=true
LockPersonality=true
SystemCallArchitectures=native
RestrictAddressFamilies=AF_INET AF_INET6 AF_NETLINK
CapabilityBoundingSet=
AmbientCapabilities=
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
